# dockerfile/download/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[program:get_gfs]
command=/bin/bash /data/collect/scripts/get_gfs.sh
